# service/asus-ryzen-power.service
[Unit]
Description=Set my custom power limits.
[Service]
ExecStart=/usr/lib/asus-ryzen-power/applypowerconfig.sh
[Install]
WantedBy=default.target
